# les1/crontab
* * * * * php /mnt/c/Users/Documents/practicecode/php-cli/les1/script.php >> /mnt/c/Users/Documents/log.txt